# bubble-server/src/main/resources/packer/roles/finalizer/files/supervisor_bubble_nodemanager.conf

[program:nodemanager]
stdout_logfile = /var/log/bubble/nodemanager-out.log
stderr_logfile = /var/log/bubble/nodemanager-err.log
command=/usr/sbin/bubble-nodemanager
