# automation/roles/bubble_finalizer/files/supervisor_bubble_nodemanager.conf

[program:nodemanager]
stdout_logfile = /home/bubble/logs/nodemanager-out.log
stderr_logfile = /home/bubble/logs/nodemanager-err.log
command=/usr/sbin/bubble-nodemanager
