# bubble-server/src/main/resources/packer/roles/algo/files/supervisor_algo_refresh_users_monitor.conf

[program:algo_refresh_users_monitor]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/bin/algo_refresh_users_monitor.sh
