# automation/roles/bubble/files/supervisor_refresh_bubble_ssh_keys_monitor.conf

[program:refresh_bubble_ssh_keys_monitor]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/sbin/refresh_bubble_ssh_keys_monitor.sh
