# automation/roles/algo/files/supervisor_wg_monitor_connections.conf

[program:wg_monitor_connections]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/sbin/wg_monitor_connections.sh
