# bubble-server/src/main/resources/packer/roles/bubble/files/supervisor_refresh_flex_keys_monitor.conf

[program:refresh_flex_keys_monitor]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/sbin/refresh_flex_keys_monitor.sh
