# bubble-server/src/main/resources/packer/roles/tarpit/files/supervisor_ssh_tarpit.conf

[program:ssh_tarpit]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/home/tarpit/bubble_ssh_tarpit.py
stopsignal=QUIT