# bubble-server/src/main/resources/packer/roles/tarpit/files/supervisor_http_tarpit.conf

[program:http_tarpit]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=sudo -u tarpit /home/tarpit/bubble_http_tarpit.py
stopsignal=QUIT