# bubble-server/src/main/resources/ansible/roles/mitmproxy/files/supervisor_mitm_monitor.conf

[program:mitm_monitor]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/sbin/mitm_monitor.sh
stopsignal=QUIT
