# bubble-server/src/main/resources/ansible/roles/mitmproxy/templates/supervisor_mitmproxy.conf.j2

[program:mitm{{ port }}]
stdout_logfile = /home/mitmproxy/mitm{{ port }}-out.log
stderr_logfile = /home/mitmproxy/mitm{{ port }}-err.log
command=sudo -H -u mitmproxy bash -c "/home/mitmproxy/mitmproxy/run_mitm.sh {{ port }}"
stopasgroup=true
stopsignal=QUIT
