# automation/roles/mitmproxy/files/supervisor_mitmdump_monitor.conf

[program:mitmdump_monitor]
stdout_logfile = /dev/null
stderr_logfile = /dev/null
command=/usr/local/sbin/mitmdump_monitor.sh
